# recipes-core/application/files/app.service
[Unit]
Description=GNU Application startup script for Arod

[Service]
ExecStart=/usr/bin/app

[Install]
WantedBy=multi-user.target
